# gorpm.service
[Unit]
Description=golang rpm

[Service]
Type=simple
ExecStart=/usr/bin/gorpm -config /etc/gorpm/config.json

[Install]
WantedBy=multi-user.target